# systemd/claude-auth-refresh.service
[Unit]
Description=Refresh Claude OAuth credentials for Uncypher containers
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=simple
# Replace DEPLOY_USER with the host user that ran `claude login` and owns the
# Docker Compose deployment (e.g. ec2-user, ubuntu, deploy). Replace
# /home/DEPLOY_USER and the uncypher path below to match that user's home and
# the directory containing docker-compose.yml. See claude-auth/README.md.
User=DEPLOY_USER
WorkingDirectory=/home/DEPLOY_USER/uncypher
Environment=COMPOSE_PROJECT_NAME=uncypher
Environment=COMPOSE_DIR=/home/DEPLOY_USER/uncypher
Environment=HOME=/home/DEPLOY_USER
Environment=PATH=/home/DEPLOY_USER/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/DEPLOY_USER/uncypher/scripts/claude-oauth-refresh-daemon.sh
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
